# repo/secfix/Library/LaunchDaemons/io.github.JPaulMora.SecFix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.myrepospace.ThePlayer10.SecFix</string>
	<key>ProgramArguments</key>
	<array>
		<string>bash</string>
		<string>/usr/bin/foo.bash</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>KeepAlive</key>
	<dict>
		<key>OtherJobEnabled</key>
		<dict>
			<key>com.apple.wifid</key>
			<true/>
		</dict>
		<key>NetworkState</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>AbandonProcessGroup</key>
	<true/>
</dict>
</plist>
